# etc/cron.d/network_check
*/5 * * * * root /usr/local/share/CasjaysDev/scripts/bin/am_i_online check
