# etc/cron.d/run-os-update
# Update system configurations everyday at 4am
0 4 * * * root [ -x "/root/.local/bin/run-os-update" ] && /root/.local/bin/run-os-update scripts >/dev/null 2>&1 && /root/.local/bin/run-os-update >/dev/null 2>&1 || bash -c "$(curl -q -LSsf "https://github.com/casjay-base/centos/raw/main/root/.local/bin/run-os-update")" >/dev/null 2>&1
